# core/openclaw/native/systemd/openclaw.service.example
[Unit]
Description=fusionAIze OpenClaw (native)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Group=openclaw
Environment=HOME=/var/lib/openclaw
WorkingDirectory=/var/lib/openclaw

# pin profile
Environment=OPENCLAW_PROFILE=prod

# env files (created on core by install.sh)
EnvironmentFile=/etc/openclaw/openclaw.env
EnvironmentFile=/etc/openclaw/openclaw.token.env

ExecStart=/usr/local/bin/openclaw gateway run \
  --profile prod \
  --allow-unconfigured \
  --force --verbose --ws-log compact \
  --bind ${OPENCLAW_GATEWAY_BIND} \
  --port ${OPENCLAW_GATEWAY_PORT} \
  --auth ${OPENCLAW_GATEWAY_AUTH} \
  --token ${OPENCLAW_GATEWAY_TOKEN}

Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
